# code_to_copy/backend/sa/crontab
00 6    * * *   root    ntpdate 210.72.145.44
